# robot-tools/etc/cron.d/data_folder
*/NEW_DATA_FOLDER_TIME * * * * TEMPLATE_USER /bin/bash /opt/bot/scripts/update_data_directory.sh >> /tmp/z_data_folder.log
@reboot TEMPLATE_USER /bin/bash /opt/bot/scripts/update_data_directory.sh >> /tmp/z_data_folder.log
