# robot-tools/etc/cron.d/delete_data
4-59/5 * * * * TEMPLATE_USER /bin/bash /opt/bot/scripts/delete_old_files.sh MOUNT_DRIVE FOLDER PERCENT >> /tmp/z_delete_data.log
