# TEMPLATE_PROCESS/etc/cron.d/keep_running_TEMPLATE_PROCESS
* * * * * TEMPLATE_USER /bin/bash /opt/TEMPLATE_PROCESS/scripts/keepRunning.sh